# supervisord.conf.template
[supervisord]
nodaemon=true

[program:searxng]
command=/bin/bash -lc "exec /opt/searxng/searx-pyenv/bin/python searx/webapp.py"
directory=/opt/searxng/searxng-src
environment=SEARXNG_SETTINGS_PATH="/etc/searxng/settings.yml"
autostart=true
autorestart=true
startsecs=15
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pm2]
command=/bin/bash -lc "source /usr/local/nvm/nvm.sh && nvm use 12 && export OPENAI_API_KEY=\"${OPENAI_API_KEY:-}\" && pm2 start /root/c9sdk/server.js --name=\"c9sdk\" -- -w / -l 0.0.0.0 -p 3399 -a c9sdk:%(ENV_C9SDK_PASSWORD)s && nvm use 24 && pm2 start /root/.openclaw/pm2.config.js --only openclaw-gateway && pm2 start /root/auto-approve-openclaw/pm2.config.js --only openclaw-auto-approve && pm2 start /var/www/html/ide-readme/pm2.config.js --only readme"
priority=4

[program:install-included-skills]
command=/bin/bash -lc "source /usr/local/nvm/nvm.sh && nvm use 24 >/dev/null && exec /var/www/html/scripts/install-included-skills.sh"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/bin/bash -lc "rm -f /etc/nginx/sites-enabled/*.template 2>/dev/null || true; nginx -s quit 2>/dev/null || true; nginx -g 'daemon off;'"
priority=20

[program:cron]
command=cron -f

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +extension RANDR +render -noreset
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce]
command=/usr/local/bin/start-xfce-session.sh
environment=DISPLAY=":1",HOME="/root",USER="root",XDG_RUNTIME_DIR="/tmp/runtime-root"
autostart=true
autorestart=true
startsecs=5
priority=6
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/local/bin/start-x11vnc.sh
environment=DISPLAY=":1",VNC_PORT="5901",VNC_PASSWORD="%(ENV_VNC_PASSWORD)s"
autostart=true
autorestart=true
startsecs=3
priority=7
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/bin/websockify --web=/opt/novnc 6080 127.0.0.1:5901
autostart=true
autorestart=true
startsecs=3
priority=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0